# examples/systemd/salt-security-agent-beat.service
[Unit]
Description=Salt Security Agent - Celery Beat scheduler
Documentation=https://github.com/
After=network.target redis.service salt-security-agent-worker.service
Wants=redis.service

[Service]
Type=simple
User=salt-security-agent
Group=salt-security-agent
WorkingDirectory=/etc/salt-security-agent
ExecStart=/opt/salt-security-agent/venv/bin/salt-security-agent beat --config /etc/salt-security-agent/config.toml --loglevel INFO
Restart=on-failure
RestartSec=5

KillMode=mixed
TimeoutStopSec=30

# Hardening. Beat writes its schedule database; by default Celery places it
# in the working directory (celerybeat-schedule). Adjust ReadWritePaths if
# you point Beat at a different location.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/salt-security-agent

[Install]
WantedBy=multi-user.target
