# examples/systemd/salt-security-agent-worker.service
[Unit]
Description=Salt Security Agent - Celery worker
Documentation=https://github.com/
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=salt-security-agent
Group=salt-security-agent
WorkingDirectory=/etc/salt-security-agent
ExecStart=/opt/salt-security-agent/venv/bin/salt-security-agent worker --config /etc/salt-security-agent/config.toml --loglevel INFO
Restart=on-failure
RestartSec=5

# Give Celery time to finish in-flight scans before being killed.
KillMode=mixed
TimeoutStopSec=30

# Hardening. The salt CLI (salt, salt-key) and the Salt state repo
# (config: [salt].repo_path) must be reachable for the salt-security-agent user.
# Extend ReadWritePaths / ReadOnlyPaths as needed for your environment.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/salt-security-agent

[Install]
WantedBy=multi-user.target
